# tpl/e46b5789ce744577a2f91b2f06d50fd6/etc/systemd/system/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty --autologin dimitri --noclear %I 38400 linux
