# tpl/defaults/etc/systemd/user/udiskie.service
[Unit]
Description=Automounting daemon

[Service]
ExecStart=/usr/bin/udiskie
Restart=always

[Install]
WantedBy=default.target

