# tpl/e46b5789ce744577a2f91b2f06d50fd6/etc/systemd/system/gpg-luks-mount@.service
[Unit]
Description=Mount gpg protected luks devices

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gpg-luks.sh -k /etc/luks-1234567AD353.gpg -d %I
