# tpl/e46b5789ce744577a2f91b2f06d50fd6/etc/systemd/system/gpg-luks-unmount@.service
[Unit]
Description=Unmount gpg protected luks devices

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gpg-luks.sh -u -d %I
